# gamebbs.conf
[program:gamebbs]
command=supervisor -i ./public bin/www
directory=/var/www/bbs
autostart=true
autorestart=true
redirect_stderr=true